# pkgroot/Library/LaunchDaemons/io.micromdm.depnag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.micromdm.depnag</string>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>ThrottleInterval</key>
        <integer>60</integer>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/depnag</string>
        </array>
        <key>StandardErrorPath</key>
        <string>/var/log/micromdm/depnag-stderr.log</string>
        <key>StandardOutPath</key>
        <string>/var/log/micromdm/depnag-stdout.log</string>
    </dict>
</plist>

